# supervisor/boundaryissues-supervisor.conf.example
[program:boundaryissues_gearman_worker]
command=/usr/bin/php -f /usr/local/mapzen/whosonfirst-www-boundaryissues/bin/gearman_worker.php
user=www-data
